# non_stow/root/etc/systemd/system/borg.service
[Unit]
Description=Automated Borg Backup:

[Service]
Type=oneshot


ExecStart=bash /home/oskar/borg.sh
